# deploy/native/openclaw-proxy.service
[Unit]
Description=OpenClaw Secure Proxy
After=network-online.target openclaw.service
Wants=network-online.target
Requires=openclaw.service
Documentation=https://github.com/yihuang/openclaw-secure-stack

[Service]
Type=simple
User=ocproxy
Group=ocproxy
WorkingDirectory=/opt/openclaw-secure-stack
EnvironmentFile=/etc/openclaw-secure-stack/proxy.env

# Wait for OpenClaw to be healthy before starting (30-second timeout)
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do curl -sf http://127.0.0.1:3000/health && exit 0 || sleep 2; done; exit 1'

# Start uvicorn with 2 workers (optimal for 4-thread CPU)
ExecStart=/opt/openclaw-secure-stack/.venv/bin/python -m uvicorn \
    src.proxy.app:create_app_from_env \
    --host 127.0.0.1 --port 8080 --factory --workers 2

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# Security: systemd sandboxing (replaces Docker read_only, cap_drop, no-new-privileges)
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictSUIDSGID=yes

# Python: pydantic_core uses precompiled binaries, test if MemoryDenyWriteExecute=yes works
# If startup fails with "cannot map memory", change to =no
MemoryDenyWriteExecute=yes

# Filesystem permissions: RW for DBs and audit logs, RO for config
ReadWritePaths=/var/lib/openclaw-proxy ${HDD_MOUNT}/openclaw-audit
ReadOnlyPaths=/opt/openclaw-secure-stack/config

# Resource limits (2010 Mac Mini: 8GB RAM, 4 threads)
MemoryMax=2G
TasksMax=128

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-proxy

[Install]
WantedBy=multi-user.target
